# dist/orangepi-release/deploy/orangepi/piso-wifi.service
[Unit]
Description=Piso WiFi Portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/piso-wifi
ExecStart=/opt/piso-wifi/deploy/orangepi/run.sh
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
